# deploy/production/templates/finance-api.service
[Unit]
Description=Finance Go API
Documentation=file:/opt/finance/current/operations/production-runbook.md
After=network-online.target
Wants=network-online.target
ConditionPathIsExecutable=/opt/finance/current/bin/finance-api
StartLimitIntervalSec=60s
StartLimitBurst=5

[Service]
Type=simple
User=finance
Group=finance
WorkingDirectory=/opt/finance/current
EnvironmentFile=/etc/finance/finance-api.env
ExecStart=/opt/finance/current/bin/finance-api
Restart=on-failure
RestartSec=5s
TimeoutStartSec=30s
TimeoutStopSec=90s
KillSignal=SIGTERM
KillMode=mixed
UMask=0077

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectClock=true
ProtectHostname=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RemoveIPC=true
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
ReadOnlyPaths=/opt/finance /etc/finance

LimitNOFILE=4096
TasksMax=128
MemoryMax=768M
CPUQuota=200%

StandardOutput=journal
StandardError=journal
SyslogIdentifier=finance-api

[Install]
WantedBy=multi-user.target
